# s6/simple-with-env/run
#!/command/execlineb -P
# s6-setuidgid root
s6-envdir /env /app/s6-overlay-showcase
